# cookbooks/basics/templates/default/rc.local.erb
#!/bin/sh
#
# Generated by Chef for <%= node[:fqdn] %>
# Local modifications will be overwritten.
export LANG=en_US.UTF-8
touch /var/lock/subsys/local
sed -i -r 's/HOSTNAME.*/HOSTNAME=<%= node[:fqdn] %>/g' /etc/sysconfig/network
hostname <%= node[:fqdn] %>

service ntpd stop
ntpdate -b 50.116.38.157
service ntpd start
